# app/seaking_server/crontab
* * * * * su seaking -c /home/wozlla/redis-2.6.14/nodeprocess.sh